# supervisor.flask.conf
[program:flask]
command = gunicorn app:app -b 0.0.0.0:5000
directory = /app
autostart = true
autorestart = true