# supervisor.py3o-fusion.conf
[program:renderserver]
command = start-py3o-fusion --renderserver=localhost
autostart = true
autorestart = true